# usr/lib/systemd/system/livecheck-lsblk.service
[Unit]
Description=Obtains lsblk output for use by livecheck

[Service]
Type=oneshot
RemainAfterExit=no
ExecStartPre=bash -c 'safe-rm -r -f -- /run/desktop-config-dist'
ExecStart=/usr/share/livecheck/livecheck-lsblk
ExecStartPost=bash -c 'touch /run/desktop-config-dist/done'

[Install]
WantedBy=multi-user.target
