# etc/xdg/autostart/backlight-tool-dist-restore.desktop
[Desktop Entry]
Type=Application
Name=Backlight state restore
Exec=/usr/bin/backlight-tool-dist-agent restore
StartupNotify=false
NoDisplay=true
