# etc/xdg/autostart/livecheck.desktop
[Desktop Entry]
Type=Application
Name=Live status checker
Exec=systemctl --user start livecheck-systray.service
StartupNotify=false
NoDisplay=true
