# etc/xdg/autostart/livecheck.desktop
[Desktop Entry]
Type=Application
Name=Live status checker
Exec=livecheck-systray
StartupNotify=false
NoDisplay=true
